# examples/99-backlight.udev.rules
SUBSYSTEM=="backlight", ACTION=="add", \
  RUN+="/bin/chgrp wheel %S%p/brightness", \
  RUN+="/bin/chmod g+w %S%p/brightness"
